# borg/.config/systemd/user/bbkp.timer
[Unit]
Description=backup user files every 4 hours of user uptime

[Timer]
OnUnitActiveSec=4h

[Install]
WantedBy=timers.target
