# private_dot_config/systemd/user/bbkp.service
[Unit]
Description=run borg backup

[Service]
ExecStart=bbkp

[Install]
WantedBy=default.target
